# supervisor_celery.conf
[program:celery]
command = celery -A stackviz_deployer.tasks.tasks worker
directory = /home/docker/stackviz-deployer
user = celery
